# data/indicator-sound.conf.in
description "Indicator Sound Backend"

# Want to move to indicator-services-[start|end], but that's not all
# there yet.  Use the signals that exist today for now.

start on indicators-loaded or indicator-services-start
stop on desktop-end or indicator-services-end

respawn

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/indicator-sound/indicator-sound-service
