# data/indicator-sound.service.in
[Unit]
Description=Indicator Sound Service
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/indicator-sound/indicator-sound-service
Restart=on-failure
